# deploy/systemd/asset-master-data.slice.tpl
[Unit]
Description=Aggregate resource budget for asset master-data

[Slice]
MemoryAccounting=true
MemoryHigh=224M
MemoryMax=288M
TasksMax=192
